# backend/supervisor_blog_cron.conf
[program:blog-cron]
command=/usr/bin/python3 /app/backend/blog_cron.py
directory=/app/backend
autostart=false
autorestart=false
startsecs=0
stdout_logfile=/var/log/supervisor/blog-cron.out.log
stderr_logfile=/var/log/supervisor/blog-cron.err.log
